# package/daemon.rpm.prerm.in
@DAEMON_PRERM_IN_BODY@
chkconfig --del lizardfs-@COMPONENT_NAME@

