# package/daemon.rpm.postinst.in
@DAEMON_POSTINST_IN_BODY@
chkconfig --add lizardfs-@COMPONENT_NAME@

