# package/daemon.deb.postinst.in
@DAEMON_POSTINST_IN_BODY@
update-rc.d lizardfs-@COMPONENT_NAME@ defaults
